# src/main/deb/etc/systemd/system/nginj.service
[Unit]
Description="desc"
[Service]
ExecStart=/usr/bin/java -jar /usr/lib/nginj/nginj.jar -Xmx8G -Xms8G -server -Djava.awt.headless=true
Type=simple